# ServerB/cron.d/17off
31 14 * * * root /home/pi/gpio17off.sh 2>&1 ; mail -e -s Apagado proyectosoii2020@gmail.com
