# ServerB/cron.d/17on
  * * * root /home/pi/gpio17on.sh 2>&1 ; mail -e -s Encendido proyectosoii2020@gmail.com
