# ServerA/cron.d/27on
06 18 * * * root /home/pi/gpio27on.sh 2>&1 ; mail -e -s Encendido proyectosoii2020@gmail.com
